# deploy/protectinator-fleet-scan.timer
[Unit]
Description=Run Protectinator fleet scan every 6 hours

[Timer]
OnCalendar=*-*-* 00/6:00:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
